# home/private_dot_config/systemd/user/orca-server.service.tmpl
{{- if and (eq .chezmoi.os "linux") .orcaServer -}}
{{- /* CODER_WORKSPACE_NAME is the canonical workspace identifier on Coder
       boxes, set by the Coder agent in /run/coder/agent.conf. We extract it
       at chezmoi-render time so the systemd unit gets a baked-in pairing
       address — no runtime shell parsing, no risk of leaking sibling secrets
       from agent.conf (e.g. CODER_AGENT_TOKEN) into the orca process. */ -}}
{{- $agentConf := "" -}}
{{- if stat "/run/coder/agent.conf" -}}
{{-   $agentConf = include "/run/coder/agent.conf" -}}
{{- end -}}
{{- $workspaceName := $agentConf | regexFind `CODER_WORKSPACE_NAME="[^"]+"` | regexFind `[^"]+"$` | trimSuffix `"` -}}
{{- if eq "" $workspaceName -}}
{{-   fail "orcaServer = true but CODER_WORKSPACE_NAME could not be resolved from /run/coder/agent.conf" -}}
{{- end -}}
[Unit]
Description=Orca server (headless runtime for remote pairing)
After=network.target

[Service]
ExecStart={{ .chezmoi.homeDir }}/.local/bin/orca-serve --pairing-address {{ $workspaceName }}.coder --json
Environment=HOME={{ .chezmoi.homeDir }}
Environment=PATH={{ .chezmoi.homeDir }}/.local/bin:/usr/bin:/usr/local/bin:/bin
SuccessExitStatus=1 143
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
{{- end }}
